# xbow/Packer/xbow-worker.path
[Unit]
Description='Watch for creation of /run/metadata/xbow/scheduler_ip_address'

[Path]
PathExists=/run/metadata/xbow/scheduler_ip_address

[Install]
WantedBy=multi-user.target
